# systemd/.config/systemd/user/hyprland.service
[Unit]
Description=A service for hyprland wrapper

[Service]
Type=simple
ExecStart=%h/.local/bin/wrappedhl

[Install]
WantedBy=default.target
